# packaging/powertooth.service
[Unit]
Description=PowerTooth BlueZ gamepad bridge
After=bluetooth.service dev-powertooth.device
Wants=bluetooth.service

[Service]
Type=simple
ExecStart=/usr/local/bin/powertooth --device /dev/powertooth
Restart=on-failure
RestartSec=2
SupplementaryGroups=dialout
StandardOutput=append:/var/log/powertooth/host.log
StandardError=append:/var/log/powertooth/host.log
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
